# systemd/burp-server.socket
[Unit]
Description=Burp backup server socket listener

[Socket]
ListenStream=4971

[Install]
WantedBy=sockets.target
